# nginx/rootfs/etc/s6-overlay/s6-rc.d/nginx-prepare/up
#!/command/execlineb -P
elglob -0 FILES /var/tmp/nginx/*
rm -rf -- $(FILES)
