# email-oauth2-proxy/rootfs/etc/s6-overlay/s6-rc.d/emailproxy-prepare/up
/etc/s6-overlay/scripts/emailproxy-prepare
